# papertrail.service
[Unit]
Description=Papertrail GPS Tracker
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/papertrail
ExecStart=/home/pi/papertrail/start.sh
Restart=on-failure
RestartSec=10
StandardOutput=append:/var/log/papertrail.log
StandardError=append:/var/log/papertrail.log

# Environment
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
